# build/deb/stable/etc/cron.d/open-semantic-search
#
# cron-jobs for Open Semantic Search Apps
#

# RSS Feeds import
*/1 * * * *     opensemanticetl curl --silent http://localhost/search-apps/rss_manager/import

# Webcrawler
*/1 * * * *     opensemanticetl curl --silent http://localhost/search-apps/crawler/recrawl

# File crawler
*/1 * * * *     opensemanticetl curl --silent http://localhost/search-apps/files/recrawl
